# LaunchAgent/com.tandemclip.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tandemclip</string>

    <!-- Default: the .app built by Scripts/make-app.sh and copied to
         /Applications. For a bare-binary install use /usr/local/bin/tandemclip.
         Add "--verbose" as a second <string> to trace discovery/TLS/sync. -->
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/TandemClip.app/Contents/MacOS/tandemclip</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <!-- No StandardOutPath/StandardErrorPath: writing to world-readable /tmp
         leaked logs to other local users. The app logs via NSLog to the unified
         logging system (access-controlled); enable "verbose" in Settings and
         read it with Console.app or:
         log stream --predicate 'process == "tandemclip"' -->
</dict>
</plist>
